# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:webman]
command=php /app/start.php start
directory=/app
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/webman.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:cron]
command=cron -f
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cron.log

[program:auto-deploy]
command=/usr/local/bin/auto-deploy.sh
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/auto-deploy.log
